# .config/cron/crontab
0 * * * * /Users/kenny/.config/scripts/news-update.sh >> /Users/kenny/.config/cron/logs/news-update.log 2>&1
